# vm/80-qubes-backlight.rules
ACTION=="change", KERNEL=="dummy_backlight", SUBSYSTEM=="backlight", RUN+="/usr/lib/qubes/qubes-set-backlight.sh"
